# .config/yadm/system_alt/etc/systemd/user/unison@.service
[Unit]
Description= Unison Profile Watcher

[Service]
Type=simple
ExecStart=/usr/bin/unison %i

[Install]
WantedBy=default.target
