# .config/yadm/system_alt/etc/cron.monthly/btrfs-scrub##c.openrc
#!/bin/sh

/usr/bin/btrfs scrub start -B / >> /var/log/btrfs-scub.log 2>&1
